# camera/client/roles/cam_view/files/getty_tty1_override.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin cam_view --noclear %I 38400 linux
